# src/Lamashare.CLI/Lamashare.CLI/bin/Debug/net8.0/Resources/lamashare.service
[Unit]
Description=Lamashare background services

[Service]
ExecStart={ExecStart}
Restart=always

[Install]
WantedBy=default.target